# deploy/systemd/kekkai-agent.service
[Unit]
Description=kekkai edge XDP firewall agent
Documentation=https://github.com/ExpTechTW/kekkai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kekkai-agent -config /etc/kekkai/kekkai.yaml -managed-config /etc/kekkai/kekkai.agent.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2s

# Capabilities — run as root (needed for bpffs + XDP attach on most kernels)
# but drop every capability except what XDP actually needs. CAP_SYS_ADMIN is
# required because cilium/ebpf uses bpf() syscall paths that check for it on
# kernels where CAP_BPF is not sufficient.
User=root
AmbientCapabilities=CAP_BPF CAP_NET_ADMIN CAP_PERFMON CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_BPF CAP_NET_ADMIN CAP_PERFMON CAP_SYS_ADMIN

# Hardening — keep the agent's blast radius small.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false
# Writable paths the agent actually needs.
ReadWritePaths=/sys/fs/bpf /var/run /run /etc/kekkai

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kekkai-agent

[Install]
WantedBy=multi-user.target
